# tartan/basic test/pinere-boot-test-service
#!/sbin/openrc-run
# PineRE Tartan Early Boot Service

name="PineRE Boot Test"
description="Display selected boot mode and verify system startup"

depend() {
    need localmount
}

start() {
    ebegin "Starting PineRE boot test"

    BOOTMODE="es"

    if grep -q "bootmode=shell" /proc/cmdline; then
        BOOTMODE="shell"
    elif grep -q "bootmode=xfce" /proc/cmdline; then
        BOOTMODE="xfce"
    elif grep -q "bootmode=waydroid" /proc/cmdline; then
        BOOTMODE="waydroid"
    elif grep -q "bootmode=menu" /proc/cmdline; then
        BOOTMODE="menu"
    fi

    /usr/bin/pinere-boot-test "${BOOTMODE}"

    eend $?
}
